# systemd/user/firefox@Smaato.service.d/delay.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/bin/sleep 5
