# autostart/feeds.desktop
[Desktop Entry]
Comment=RSS reader
Exec=gfeeds
Name=Feeds
Terminal=false
Type=Application
