# autostart/firefox-Smaato.desktop
[Desktop Entry]
Comment=Firefox with Smaato profile
Exec=firefox -P Smaato
Name=Firefox
Terminal=false
Type=Application
X-GNOME-Autostart-Delay=5
