# autostart/lollypop.desktop
[Desktop Entry]
Comment=Music player
Exec=lollypop
Name=Lollypop
Terminal=false
Type=Application
