# autostart/terminal.desktop
[Desktop Entry]
Exec=gnome-terminal --wait --maximize
Name=Terminal
Terminal=false
Type=Application
