# systemd/user/keyboxd.socket
[Unit]
Description=GnuPG key database daemon

[Socket]
ListenStream=%t/gnupg/S.keyboxd
