# systemd/user/bash-completion@npm.service.d/settings.conf
[Service]
#ExecStart=/usr/bin/npm completion
ExecStart=/usr/bin/echo foo bar
